# bokeh/templates/bokeh-server.conf
{%- from 'nginx/settings.sls' import NUM_SERVERS, SERVER_NAME, PUBLIC_IP with context -%}
{%- from 'bokeh/settings.sls' import BOKEH_DEMOS_LIST with context -%}

{%- set environment = [] -%}
{%- do environment.append('HOME="/home/ec2-user"') -%}
{%- do environment.append('USER="ec2-user"') -%}
{%- do environment.append('PATH="/home/ec2-user/miniconda/envs/bokeh/bin:%(ENV_PATH)s"') -%}

[program:bokeh_demos]
command=/home/ec2-user/miniconda/envs/bokeh/bin/bokeh serve {{ BOKEH_DEMOS_LIST | join(' ') }} --prefix=apps/ --address=127.0.0.1 --port 51%(process_num)02d --host {{SERVER_NAME}}:443 --host {{PUBLIC_IP}}:443 --allow-websocket-origin={{SERVER_NAME}} --use-xheaders
directory=/home/ec2-user/
autostart=false
autorestart=true
startretries=3
numprocs={{ NUM_SERVERS }}
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/home/ec2-user/log/bokehserver.err.log
stdout_logfile=/home/ec2-user/log/bokehserver.out.log
user=ec2-user
environment={{ environment | join(',') }}
