# bokeh/templates/bokeh-server.conf
{%- from 'nginx/settings.sls' import NUM_SERVERS, PUBLIC_IP with context -%}
{%- from 'bokeh/settings.sls' import BOKEH_DEMOS_LIST with context -%}

[program:bokeh_demos]
command=/home/ec2-user/miniconda/envs/bokeh/bin/bokeh serve {{ BOKEH_DEMOS_LIST | join(' ') }} --prefix=apps/ --address=127.0.0.1 --port 51%(process_num)02d --host {{PUBLIC_IP}}:8080
directory=/home/ec2-user/
autostart=false
autorestart=true
startretries=3
numprocs={{ NUM_SERVERS }}
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/home/ec2-user/log/sliders.err.log
stdout_logfile=/home/ec2-user/log/sliders.out.log
user=ec2-user
environment=USER="ec2-user",HOME="/home/ec2-user"
